# _linux/var/spool/cron/crontabs/root
* * * * * echo "Hello world" >> /var/log/cron.log 2>&1
* * * * * /usr/local/bin/php /app/artisan schedule:run 1>> /var/log/cron.log 2>&1
* * * * * env >> /var/log/cron.log  >> /var/log/cron.log 2>&1
